# rpi/systemd/gpio-rcv.service
[Unit]
Description=GPIO Rcv
After=network-online.target

[Service]
ExecStart=/home/showtimer/src/showtimer/rpi/gpio-rcv.py
#KillMode=process
#Restart=on-failure
Type=simple
User=showtimer
WorkingDirectory=/home/showtimer/src/showtimer/rpi

[Install]
WantedBy=multi-user.target
